# deploy/systemd/goobreview.service.example
[Unit]
Description=Run GoobReview PR reviewer once
Documentation=https://github.com/asavs/goobreview
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=goobreview
Group=goobreview
WorkingDirectory=/opt/goobreview/example
Environment=REVIEWER_ENV_FILE=/opt/goobreview/example/config/reviewer.env
ExecStartPre=/usr/bin/bash scripts/reviewer/require-dry-run-artifact.sh
ExecStart=/usr/bin/bash scripts/reviewer/run-once.sh
